# call-forwarding-monitor.timer
[Unit]
Description=Run Call Forwarding Health Check every 5 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
Persistent=true

[Install]
WantedBy=timers.target
